# etc/init/shutdown.conf
start on runlevel RUNLEVEL=0

kill timeout 3600

task

emits shutdown

script 
  initctl emit shutdown
end script
