# etc/init/yavdr-backup.conf
start on shutdown

task

console output

kill timeout 3600

script
   /usr/bin/yavdr-backup
end script
